# cron.d/certbot
30 */12 * * * root /opt/housekeeping/bin/certbot-renew
